# systemd/fire-kiosk.service
[Unit]
Description=Fireplace Kiosk Watcher
After=multi-user.target fire-web.service
Wants=fire-web.service

[Service]
Type=simple
# User will be set during installation based on actual system user
# Default: User=fireplace, but should be changed to actual user during setup
User=fireplace
Group=fireplace
WorkingDirectory=/opt/fireplace
Environment="PATH=/opt/fireplace/venv/bin:/usr/bin:/bin"
Environment="DISPLAY=:0"
# These paths will be set during installation based on actual user
# Environment="XAUTHORITY=/home/[USER]/.Xauthority"
# Environment="HOME=/home/[USER]"
# For now, let the watcher.py auto-detect these paths

# Use the launch script which handles environment setup
ExecStart=/opt/fireplace/scripts/launch_kiosk.sh

# Restart configuration
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryMax=1G
CPUQuota=80%

[Install]
WantedBy=default.target